# ansible/mysql-mss-mha/templates/crontab
0 4 * * * /bin/bash /root/purge_relay_log.sh